# install/etc/cron.d/quotek
* * * * * root php /usr/local/quotek/admin/jobs/scheduler.php 2>&1 >> /var/log/quotek-admin-scheduler.log
